# server/distribution/all/conf/supervisord/registry.supervisord.conf
[program:registry]
command=/registry-distribution/registry-all/bin/registry-run.sh
startsecs=30
autostart = true
autorestart=true